# roles/backup/templates/etc_systemd_system_rsnapshot-daily.service.j2
[Unit]
Description=rsnapshot daily backup

[Service]
Type=oneshot
ExecStart=/usr/bin/rsnapshot daily
SyslogIdentifier=rsnapshot-daily

[Install]
WantedBy=multi-user.target
