# roles/backup/templates/etc_systemd_system_rsnapshot-monthly.service.j2
[Unit]
Description=rsnapshot monthly backup

[Service]
Type=oneshot
ExecStart=/usr/bin/rsnapshot monthly
SyslogIdentifier=rsnapshot-monthly

[Install]
WantedBy=multi-user.target
