# roles/backup/templates/etc_systemd_system_rsnapshot-daily.timer.j2
[Unit]
Description=rsnapshot daily backup

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*-*-* 01:00:00
Unit=rsnapshot-daily.service
