# roles/backup/templates/etc_systemd_system_rsnapshot-weekly.timer.j2
[Unit]
Description=rsnapshot weekly backup

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=sun 00:30:00
Unit=rsnapshot-weekly.service
