# roles/backup/templates/etc_systemd_system_rsnapshot-weekly.service.j2
[Unit]
Description=rsnapshot weekly backup

[Service]
Type=oneshot
ExecStart=/usr/bin/rsnapshot weekly
SyslogIdentifier=rsnapshot-weekly

[Install]
WantedBy=multi-user.target
