# roles/backup/templates/etc_systemd_system_rsnapshot-monthly.timer.j2
[Unit]
Description=rsnapshot monthly backup

[Install]
WantedBy=timers.target

[Timer]
OnCalendar=*-*-01 04:01:00
Unit=rsnapshot-monthly.service
