# roles/monitoring_goaccess/templates/etc_systemd_system_goaccess-update.timer.j2
[Unit]
Description=Generate goaccess HTML report

[Install]
WantedBy=timers.target

[Timer]
Unit=goaccess-update.service
Persistent=true
OnCalendar={{ goaccess_update_calendar_expression }}
